# crons/root.dev7.cron


*/5 * * * * root /bin/sh /opt/vyew/bin/chkCpuHog dev7.vyew.com

